# seedbox-web.service
[Unit]
Description=Seedbox Web Dashboard
After=network.target

[Service]
Type=simple
User=amit
WorkingDirectory=/home/amit/seedbox-web
ExecStart=/home/amit/seedbox-web/venv/bin/python app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
